# AI OS Trixie Build/payload/etc/systemd/system/cogos-wine-bridge.service
[Unit]
Description=Wolf CoG OS UL App Bridge (wine-wolf-bridge)
After=network.target

[Service]
Type=simple
Environment=COGOS_ROOT=/opt/cogos
ExecStart=/usr/bin/python3 /opt/cogos/bin/cogos-wine-bridge.py daemon
Restart=on-failure

[Install]
WantedBy=multi-user.target
